# service/bililive.service
# /etc/systemd/system/bililive.service

[Unit]
Description=BILILIVE
Documentation=https://github.com/isdut/BiliLive
After=network.target sound.target pulseaudio.service

[Service]
ExecStart=/home/ubuntu/miniconda3/bin/python /server/BiliLive/bililive.py
Restart=always
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target